# configuration/supervisord.conf
[supervisord]

[program:agube-consumer]
command=python3 /availa-agube/agube_consumer.py
redirect_stderr=true
