# configuration/supervisord.conf
[supervisord]
logfile=/availa-agube/logs/supervisord.log
pidfile=/availa-agube/logs/supervisord.pid

[program:agube-consumer]
command=python3 /availa-agube/agube_consumer.py
redirect_stderr=true
